# stage-RQB2/08-ab-boot-support/files/systemd/rasqberry-health-check.service
[Unit]
Description=RasQberry A/B Boot Health Check
Documentation=https://github.com/JanLahmann/RasQberry-Two
After=network-online.target
Wants=network-online.target
Before=rasqberry-update-poller.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rq_health_check.py
TimeoutSec=600
Restart=no
StandardOutput=journal
StandardError=journal

# Run as root (needs to confirm boot slot)
User=root

[Install]
WantedBy=multi-user.target
